# debian/scrapyd.upstart
# Scrapy service

start on runlevel [2345]
stop on runlevel [06]

script
    [ -r /etc/default/scrapyd ] && . /etc/default/scrapyd
    logdir=/var/log/scrapyd
    exec twistd -ny /usr/share/scrapyd/scrapyd.tac \
                -u scrapy -g nogroup \
                --pidfile /var/run/scrapyd/scrapy.pid \
                -l $logdir/scrapyd.log >$logdir/scrapyd.out 2>$logdir/scrapyd.err
end script
