# debian/scrapy-service.upstart
# Scrapy service

start on runlevel [2345]
stop on runlevel [06]

script
    [ -r /etc/default/scrapy-service ] && . /etc/default/scrapy-service
    exec twistd -ny /usr/share/scrapy/scrapy.tac \
                -u scrapy -g nogroup \
                -l /var/log/scrapy/service.log \
                --pidfile /var/run/scrapy/scrapy.pid
end script
